# host-services/dsmr_bridge.service
[Unit]
Description=Minyad DSMR Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minyad
WorkingDirectory=/opt/minyad/host-services
EnvironmentFile=/opt/minyad/host-services/.env
Environment=METRICS_ADDR=192.168.110.2
ExecStart=/opt/minyad/host-services/venv/bin/python dsmr_bridge.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
